# templates/job.conf.erb
# <%= @title %>
<% if @description %>description "<%= description %>"<% end %>

start on runlevel <%= @start_runlevel %>
stop on runlevel <%= @stop_runlevel %>

# daemon or not?
expect <%= @expect %> 
respawn

<% if @use_script -%>
script
  <%= @cmd %>
end script
<% else -%>
exec <%= @cmd %>
<% end -%>
